# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=make server
